# ansible/roles/jukebox/files/hearth-jukebox.service
[Unit]
Description=Hearth Jukebox Service
After=network.target

[Service]
Type=simple
User=hearth-jukebox
Group=hearth-jukebox
WorkingDirectory=/opt/hearth/jukebox/server
ExecStart=/opt/hearth/jukebox-venv/bin/uvicorn main:app --host 127.0.0.1 --port 8766
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/srv/hearth/music
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
