# ansible/roles/admin/templates/cafebox-admin.service.j2
[Unit]
Description=CafeBox admin backend
Documentation=https://github.com/trydydd/cafe-box
After=network.target cafebox-first-boot.service
Wants=network.target

[Service]
Type=simple
User=cafebox-admin
Group=cafebox-admin
# shadow group membership is required so the backend can read /etc/shadow
# for PAM / spwd-based password verification.
SupplementaryGroups=shadow
WorkingDirectory=/opt/cafebox/admin/backend
EnvironmentFile=/etc/cafebox/admin.env
ExecStart=/opt/cafebox/admin/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5s
PrivateTmp=true

[Install]
WantedBy=multi-user.target
